# ansible/files/supervisor/happyman.conf
[program:happyman]
command=/opt/happyman/happyman
redirect_stderr=true
stdout_logfile=/var/log/happyman.log
directory=/opt/happyman